# grbox-panel.service
[Unit]
Description=GRBOX Panel
After=network.target
[Service]
ExecStart=/opt/grbox/grbox-panel
Restart=always
[Install]
WantedBy=multi-user.target
